# templates/systemd.service.erb
# This file is managed by Puppet. DO NOT EDIT.
# If you modify this, please also make sure to edit init.sh

[Unit]
Description=Time series data processing engine.
Documentation=https://docs.influxdata.com/kapacitor/
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
LimitNOFILE=65536
EnvironmentFile=-<%= @service_defaults %>
ExecStart=/usr/bin/kapacitord -config <%= @configuration_path %>/<%= @configuration_file %> $KAPACITOR_OPTS
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
